# example/chaind.conf
[program:chaind]
command=/usr/bin/chaind --home /etc/chaind start
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/chaind.out.log
stderr_logfile=/var/log/supervisor/chaind.err.log